# scanopy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:postgres]
command=postgres -D /data/scanopy/postgres_data
environment=POSTGRES_DB=scanopy,POSTGRES_USER=postgres,POSTGRES_PASSWORD=%(ENV_POSTGRES_PASSWORD)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.err
priority=100

[program:scanopy-daemon]
command=docker run --rm --name scanopy-daemon --network host --privileged -v /data/scanopy/daemon_config:/root/.config/daemon -v /var/run/docker.sock:/var/run/docker.sock:ro -e SCANOPY_LOG_LEVEL=%(ENV_SCANOPY_LOG_LEVEL)s -e SCANOPY_SERVER_URL=%(ENV_SCANOPY_SERVER_URL)s ghcr.io/scanopy/scanopy/daemon:latest
environment=SCANOPY_LOG_LEVEL=%(ENV_SCANOPY_LOG_LEVEL)s,SCANOPY_SERVER_URL=%(ENV_SCANOPY_SERVER_URL)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/daemon.log
stderr_logfile=/var/log/supervisor/daemon.err
priority=200
depends_on=postgres

[program:scanopy-server]
command=docker run --rm --name scanopy-server -p 60072:60072 -v /data/scanopy:/data -e SCANOPY_LOG_LEVEL=%(ENV_SCANOPY_LOG_LEVEL)s -e SCANOPY_DATABASE_URL=%(ENV_SCANOPY_DATABASE_URL)s -e SCANOPY_WEB_EXTERNAL_PATH=/app/static -e SCANOPY_PUBLIC_URL=%(ENV_SCANOPY_PUBLIC_URL)s -e SCANOPY_INTEGRATED_DAEMON_URL=%(ENV_SCANOPY_INTEGRATED_DAEMON_URL)s --add-host=host.docker.internal:host-gateway ghcr.io/scanopy/scanopy/server:latest
environment=SCANOPY_LOG_LEVEL=%(ENV_SCANOPY_LOG_LEVEL)s,SCANOPY_DATABASE_URL=%(ENV_SCANOPY_DATABASE_URL)s,SCANOPY_WEB_EXTERNAL_PATH=/app/static,SCANOPY_PUBLIC_URL=%(ENV_SCANOPY_PUBLIC_URL)s,SCANOPY_INTEGRATED_DAEMON_URL=%(ENV_SCANOPY_INTEGRATED_DAEMON_URL)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/server.log
stderr_logfile=/var/log/supervisor/server.err
priority=300
depends_on=postgres,scanopy-daemon

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
priority=400
