# gitea/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# Nginx service finish script

set -e

# Get the exit code of the service
exit_code="${1:-0}"

if [ "$exit_code" -eq 0 ]; then
    bashio::log.info "Nginx stopped normally"
    # Don't restart on normal exit
    exit 1
else
    bashio::log.warning "Nginx crashed with exit code $exit_code, restarting..."
    # Allow restart on crash (s6 will limit to 5 restarts)
    exit 0
fi
