# gitea/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# Nginx service for Home Assistant Ingress

set -e

# Wait for Nginx configuration to be ready
until [ -f /etc/nginx/conf.d/ingress.conf ]; do
    bashio::log.info "Waiting for Nginx configuration..."
    sleep 1
done

bashio::log.info "Starting Nginx for Home Assistant ingress..."

# Start Nginx
exec nginx -g 'daemon off;'
